# unitmgr.service
[Unit]
Description=Systemd unit manager

[Service]
Restart=always
ExecStart=/usr/bin/unitmgr -src /opt/units

[Install]
WantedBy=multi-user.target
